# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:uvicorn]
command=uvicorn main:app --host 0.0.0.0 --port 8000 --app-dir /app/Aether-Engine
directory=/app/Aether-Engine
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=30
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
startsecs=3
stopwaitsecs=10
stopsignal=QUIT
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
